# x11/.config/systemd/user/nm-applet.service
[Unit]
Description=nm-applet
Documentation=man:nm-applet(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=xsession.target
